# rules/70-touchscreen-ntrig.rules
KERNEL=="event*",ATTRS{name}=="N-trig DuoSense",SYMLINK+="twofingtouch",RUN+="/bin/chmod a+r /dev/twofingtouch"

